# launchd/com.claude-mem-os.bridge.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claude-mem-os.bridge</string>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>source "$HOME/.mcp_secrets" &amp;&amp; claude-mem-os-bridge --interval 60</string>
  </array>

  <key>StandardOutPath</key>
  <string>/tmp/claude-mem-os-bridge.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/claude-mem-os-bridge.err.log</string>
</dict>
</plist>

